# build/host-deployer/root/usr/lib/systemd/system/yunion-host-deployer.service
[Unit]
Description=Yunion Host Deploy Agent
Documentation=https://docs.yunion.cn
After=network.target httpd.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/yunion/bin/host-deployer --conf /etc/yunion/host.conf
WorkingDirectory=/opt/yunion/bin
KillMode=process
Restart=always
RestartSec=30
LimitNOFILE=500000
LimitNPROC=500000

[Install]
WantedBy=multi-user.target
