# build/glance/root/opt/yunion/share/glance/yunion-glance.service
[Unit]
Description=Yunion Cloud Image Service
Documentation=http://doc.yunionyun.com
After=network.target httpd.service

[Service]
Type=simple
User=yunion
Group=yunion
ExecStart=/opt/yunion/bin/glance --config /etc/yunion/glance/glance-api.conf
WorkingDirectory=/opt/yunion
KillMode=process
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
